# workstation/systemd/workspace-backup.timer
[Unit]
Description=Run periodic workstation workspace backups

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true
Unit=workspace-backup.service

[Install]
WantedBy=timers.target
